# files/etc/init.d/auto_resize
#!/bin/sh /etc/rc.common
# 自动扩展 EFI 模式下的 OpenWrt overlay 分区（安全版本）
# 适用于 x86_64 EFI 镜像（GPT），不会修改 /dev/sda1 (EFI 分区)
# 首次启动自动执行，完成后禁用自身

START=10
STOP=90

start() {
    logger -t auto_resize "Starting EFI-safe overlay resize..."

    [ -f /etc/config/resize_done ] && {
        logger -t auto_resize "Resize already done, skipping."
        exit 0
    }

    # 找到 overlay 所在设备（通常是 /dev/sda2）
    ROOT_DEV=$(findmnt /overlay -o SOURCE -n 2>/dev/null)
    [ -z "$ROOT_DEV" ] && ROOT_DEV=$(findmnt / -o SOURCE -n 2>/dev/null)
    [ -z "$ROOT_DEV" ] && {
        logger -t auto_resize "No root or overlay device found, aborting."
        exit 1
    }

    DISK_DEV=${ROOT_DEV%%[0-9]*}
    PART_NUM=$(echo $ROOT_DEV | grep -o '[0-9]*$')

    logger -t auto_resize "Detected root device: $ROOT_DEV (disk: $DISK_DEV, part: $PART_NUM)"

    # 只处理 GPT，避免 MBR 下的 EFI 混乱
    if ! parted $DISK_DEV print | grep -q "gpt"; then
        logger -t auto_resize "Non-GPT disk detected, skipping resize."
        exit 0
    fi

    # 确保 parted 和 resize2fs 可用
    for cmd in parted resize2fs; do
        command -v $cmd >/dev/null 2>&1 || {
            logger -t auto_resize "Missing command: $cmd"
            exit 1
        }
    done

    # 确保 EFI 分区（sda1）不被修改
    if [ "$PART_NUM" -le 1 ]; then
        logger -t auto_resize "Refusing to modify EFI partition."
        exit 1
    fi

    # 扩展分区
    logger -t auto_resize "Resizing partition $ROOT_DEV to full disk..."
    parted -s "$DISK_DEV" resizepart "$PART_NUM" 100% >/dev/null 2>&1
    if [ $? -ne 0 ]; then
        logger -t auto_resize "Partition resize failed!"
        exit 1
    fi

    # 扩展文件系统
    logger -t auto_resize "Resizing filesystem..."
    resize2fs "$ROOT_DEV" >/dev/null 2>&1
    if [ $? -eq 0 ]; then
        touch /etc/config/resize_done
        /etc/init.d/auto_resize disable
        logger -t auto_resize "Resize complete, will not run again."
    else
        logger -t auto_resize "resize2fs failed!"
    fi
}
